# config/autostart/model-camera.desktop
[Desktop Entry]
Type=Application
Name=Model Camera
Exec=/bin/sh /home/pi/model-camera/launcher.sh
